# confd/confd.service
[Unit]
Description=Confd

[Service]
ExecStart=/usr/bin/confd
Restart=always

[Install]
WantedBy=basic.target
